# .apps.conf
[program:geddy]
command = /usr/local/bin/geddy
